# enola/etc/systemd/system/enola-tor.service

[Unit]
Description=Enola Tor Service
Documentation=man:tor(1)
After=network-online.target
Wants=network-online.target

[Service]
# Simple es más confiable que forking sin PIDFile
Type=simple
User=debian-tor
Group=debian-tor
ExecStart=/usr/bin/tor -f /etc/tor/torrc
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
LimitNOFILE=4096
RuntimeDirectory=tor
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target


