# enola/etc/systemd/system/enola-smoke.service
[Unit]
Description=Enola Server - periodic smoke test
After=network.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/bin/bash /opt/enola/scripts/common/smoke_test.sh
StandardOutput=append:/var/log/enola-server/smoke_timer.log
StandardError=append:/var/log/enola-server/smoke_timer.log
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
