# enola/etc/systemd/system/enola-smoke.service
[Unit]
Description=Enola Server - periodic smoke test
Documentation=/opt/enola/scripts/common/smoke_test.sh
Wants=enola-smoke.timer

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/bin/bash -lc 'mkdir -p /var/log/enola-server; echo "[SMOKE TIMER] $(date "+%F %T") | start" >> /var/log/enola-server/smoke_timer.log; /opt/enola/scripts/common/smoke_test.sh >> /var/log/enola-server/smoke_timer.log 2>&1 || true; echo "[SMOKE TIMER] $(date "+%F %T") | end" >> /var/log/enola-server/smoke_timer.log'
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
